# etc/profile.d/miner.sh
export BELLMAN_CPU_UTILIZATION=0.9
export FIL_PROOFS_MAXIMIZE_CACHING=1
export FIL_PROOFS_USE_GPU_COLUMN_BUILDER=1
export FIL_PROOFS_USE_GPU_TREE_BUILDER=1
export FIL_PROOFS_MULTICORE_SDR_PRODUCERS=2
export LOTUS_CUSTOM_IP=10.10.0.60
export LOTUS_ALLOW_TRUNCATED_LOG=1
#export LD_LIBRARY_PATH=/usr/local/cuda-11.1/lib64
#export CUDA_HOME=/usr/local/cuda
